# systemd/deployr.timer
[Unit]
Description=Run Deployr Sync periodically

[Timer]
OnBootSec=3min
OnUnitActiveSec=3m

[Install]
WantedBy=timers.target
